# supervisor/pipeline.conf
[program:pipeline]
command=/pipeline --config=/data/pipeline.conf --log=/data/pipeline.log
autostart=true
autorestart=true
